# utils/initrd-tools/debian/initrd-tools.postrm
#!/bin/sh

set -e

if [ "$1" = purge ]; then
	rm -rf /etc/mkinitrd
	if [ -x /sbin/update-modules ]; then
		/sbin/update-modules
	fi
	update-rc.d initrd-tools.sh remove > /dev/null
fi

#DEBHELPER#
