# utils/initrd-tools/debian/initrd-tools.postrm
#!/bin/sh

set -e

if [ "$1" = purge ]; then
	rm -rf /etc/mkinitrd
	update-modules
	update-rc.d initrd-tools.sh remove > /dev/null
fi

#DEBHELPER#
